# pkg/skyrepo/skyrepo-1.3.32-1-riscv64/etc/systemd/system/install-skywire.service
[Unit]
	Description=install skywire service
	After=network-online.target
	Wants=network-online.target

	[Service]
	Type=simple
	ExecStart=/bin/install-skywire

	[Install]
	WantedBy=multi-user.target
	
